# crontab
* * * * * /usr/local/var/checkin/checkin-shot.py -r
